# containers/golang/conf/supervisor-app.conf
[program:app]
user=root
directory = %(ENV_SRV)s
command = go get ./... && go run main.go
autostart = true
autorestart = true
stdout_logfile = syslog
stderr_logfile = syslog
