# containers/base/conf/supervisor-reloader.conf
[program:reloader]
user=nobody
command = /usr/local/bin/reloader
autostart = true
autorestart = true
stdout_logfile = syslog
stderr_logfile = syslog
